# lesson6/task1_logchecker/lib/systemd/system/mychecker.service
[Unit]
Description=mychecker: scan log for keyword

[Service]
Type=oneshot
EnvironmentFile=/etc/sysconfig/mychecker
ExecStart=/bin/bash /usr/bin/mychecker.sh
