# lesson6/task1_logchecker/lib/systemd/system/mychecker.timer
[Unit]
Description=Check log for keyword

[Timer]
OnCalendar=*:*:0/30
AccuracySec=1sec
Persistent=true

[Install]
WantedBy=multi-user.target
